# ddns/ddns/systemd/kestrel-ddns.timer
[Unit]
Description=15min timer

[Timer]
OnBootSec=10
OnUnitActiveSec=10
AccuracySec=1ms

[Install]
WantedBy=timers.target